# raspberry-pi/door-control.service
[Unit]
Description=FizyoPark Kapi Role Kontrol Servisi
After=multi-user.target

[Service]
Type=simple
User=__APP_USER__
WorkingDirectory=/tmp
ExecStart=/usr/bin/python3 __APP_DIR__/raspberry-pi/door-control.py
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
